# .config/systemd/user/niri-session.target
[Unit]
Description=niri compositor session
BindsTo=wayland-session.target
Wants=wayland-session.target
After=wayland-session.target
